# systemd/screentime-guardian.service
[Unit]
Description=Screentime Guardian Daemon
Documentation=https://github.com/petrockblog/screentime-guardian
After=network.target dbus.service graphical.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/screentime-guardian -config /etc/screentime-guardian/config.yaml
Restart=on-failure
RestartSec=5

# Run as root to manage user sessions
User=root
Group=root

# Create state directory
StateDirectory=screentime-guardian
ConfigurationDirectory=screentime-guardian

# Security hardening
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=false

# Allow access to required paths
ReadWritePaths=/var/lib/screentime-guardian

# Required capabilities for session management
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_KILL CAP_SETUID CAP_SETGID

# Access to D-Bus
PrivateNetwork=false

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=screentime-guardian

[Install]
WantedBy=multi-user.target
